# raspbian/udev/91-tty_usb_speed.rules
ACTION=="add", SUBSYSTEM=="usb-serial", RUN+="/bin/bash -c '/root/set_tty_usb_speed.sh'"
