# raspbian/udev/91-tty_usb_speed.rules
ACTION=="add", SUBSYSTEM=="usb-serial", RUN+="/bin/bash -c '/etc/udev/rules.d/set_tty_usb_speed.sh'"
